# dnscrypt.conf
# custom upstart job definition for dnscrypt-proxy
# http://dnscrypt.org/

description "dnscrypt-proxy daemon"

start on runlevel [12345]
stop on runlevel [06]

expect fork

script
  test -f /etc/default/dnscrypt && . /etc/default/dnscrypt

  OPTIONS="--daemonize"

  if [ -n "$LOGFILE" ]; then
    OPTIONS="$OPTIONS --logfile=$LOGFILE"
  fi
  if [ -n "${RESOLVER_ADDRESS}" ]; then
    OPTIONS="$OPTIONS --resolver-address=${RESOLVER_ADDRESS}"
  fi
  if [ -n "${PROVIDER_NAME}" ]; then
    OPTIONS="$OPTIONS --provider-name=${PROVIDER_NAME}"
  fi
  if [ -n "${PROVIDER_KEY}" ]; then
    OPTIONS="$OPTIONS --provider-key=${PROVIDER_KEY}"
  fi

  echo "/usr/local/sbin/dnscrypt-proxy started with options: $OPTIONS"	
  exec /usr/local/sbin/dnscrypt-proxy $OPTIONS
end script
